# services/jellyfin/launchd/org.jellyfin.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>org.jellyfin.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/jellyfin</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/jellyfin.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/jellyfin-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/CURRENT_USER</string>
    </dict>
</dict>
</plist>
