# installer/build/payload/Applications/Rainbow/services/minecraft/launchd/io.papermc.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.papermc.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/java</string>
        <string>-Xms{{MC_MEMORY}}</string>
        <string>-Xmx{{MC_MEMORY}}</string>
        <string>-XX:+UseG1GC</string>
        <string>-XX:+ParallelRefProcEnabled</string>
        <string>-XX:MaxGCPauseMillis=200</string>
        <string>-jar</string>
        <string>{{MC_DIR}}/paper.jar</string>
        <string>--nogui</string>
    </array>

    <key>WorkingDirectory</key>
    <string>{{MC_DIR}}</string>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>{{MC_DIR}}/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>{{MC_DIR}}/logs/launchd-stderr.log</string>
</dict>
</plist>
