# installer/build/payload/Applications/Rainbow/services/network/launchd/rocks.rainbow.ddns.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>rocks.rainbow.ddns</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/opt/rainbow/services/network/ddns-update.sh</string>
    </array>

    <!-- Run every 5 minutes -->
    <key>StartInterval</key>
    <integer>300</integer>

    <!-- Also run at load (boot) -->
    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/opt/rainbow/logs/ddns.log</string>

    <key>StandardErrorPath</key>
    <string>/opt/rainbow/logs/ddns-error.log</string>
</dict>
</plist>
